# systemd/user/sumi-update-check.service
[Unit]
Description=Check for sumi updates and notify
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/sumi update-check
