# systemd/user/sumi-cleanup.service
[Unit]
Description=Sumi cleanup (clipboard trim, old files)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/sumi cleanup
