# systemd/user/wallust-watcher.service
[Unit]
Description=Wallust wallpaper watcher — re-themes on wallpaper change
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/bash -c 'inotifywait -m -e close_write "$HOME/.cache/sumi/current-wallpaper" | while read; do sumi wallpaper apply "$(cat $HOME/.cache/sumi/current-wallpaper)"; done'
Restart=on-failure
RestartSec=3

[Install]
WantedBy=hyprland-session.target
